# tts/whisper-stt.service
[Unit]
Description=Whisper STT Server for Brain Gateway
After=network.target

[Service]
Type=simple
User=labadmin
Group=labadmin
WorkingDirectory=/home/labadmin

# Use same environment as TTS server
Environment="PATH=/home/labadmin/qwen-tts-env/bin:/usr/local/bin:/usr/bin:/bin"
Environment="WHISPER_MODEL=large-v3"
Environment="WHISPER_DEVICE=cuda:0"
Environment="WHISPER_PORT=8003"

ExecStart=/home/labadmin/qwen-tts-env/bin/python /home/labadmin/stt_server.py

Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whisper-stt

[Install]
WantedBy=multi-user.target
